# config/udev/rules.d/00-usb-keyboard.rules
ACTION=="bind", SUBSYSTEM=="usb", DRIVERS=="usb", ATTRS{idVendor}=="1ea7", ATTRS{idProduct}=="0907", RUN+="/bin/sh  /usr/local/bin/keyboard-repeat"
